# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:backend]
command=uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 2 --loop uvloop
directory=/app/backend
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=LOG_TAG="web",PYTHONPATH="/app/backend"

[program:consumer]
command=python -m app.pipeline.consumer
directory=/app/backend
autostart=true
autorestart=true
startretries=5
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=LOG_TAG="consumer",PYTHONPATH="/app/backend"

[program:scheduler]
command=python -m app.run_scheduler
directory=/app/backend
autostart=true
autorestart=true
startretries=3
startsecs=3
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=LOG_TAG="scheduler",PYTHONPATH="/app/backend"

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
